# cloudflare_fast_ip.service
[Unit]
Description=Restart Proxy Management Service

[Service]
Type=oneshot
ExecStart=cmd
