# subprojects/v2ray-bench/systemd/v2ray-bench-autoswitch.service
[Unit]
Description=Bench nodes then apply best to v2ray.service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
Nice=19
IOSchedulingClass=best-effort
IOSchedulingPriority=7
Environment=PYTHONUNBUFFERED=1

# 自动测速；仅当新最快节点速度 > 当前运行节点速度时才切换
ExecStart=/usr/local/bin/v2ray-bench-autoswitch --bytes 5000000 --timeout 25 --bench-top 5 --min-delta-bps 1
